# dashboard.conf
; -w depend on the memory limit
; -worker-connections depend on the tolerance of delay
; worker num * worker connections should less than the db max connection number
[program:dashboard]
command=/opt/env/bin/gunicorn --bind 127.0.0.1:5000 -w 20 --worker-connections 3 -k gevent -t 60 --max-requests-jitter 5000 wsgi:application
process_name=dashboard_%(process_num)02d
numprocs=1
directory=/opt/Dashboard/dashboard
redirect_stderr=true
stdout_logfile=/var/log/dashboard/app.log
autostart=true
autorestart=true
